# script/supervisor.ini
[program:pocket]
command=/home/yang/workspace/ezlog3/env/bin/gunicorn runserver:app --bind 127.0.0.1:9001
directory=/home/yang/workspace/ezlog3
user=yang
autostart=true
autorestart=true